# secretary-hub.service
[Unit]
Description=Secretary Hub - LINE message dispatcher
After=network.target

[Service]
Type=simple
# User=your-user
# WorkingDirectory=/path/to/line-secretary
# EnvironmentFile=/path/to/line-secretary/.env
# ExecStart=/path/to/line-secretary/.venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 3939
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
